# scripts/etc/init.d/qmanager_vpn_zone
#!/bin/sh /etc/rc.common
# =============================================================================
# qmanager_vpn_zone — Boot-Time Tailscale Firewall Zone Self-Heal
# =============================================================================
# Runs once at boot, GATED on the user-controlled toggle
# quecmanager.tailscale_workarounds.enabled. When the flag is "1", ensures the
# fw4 zone for tailscale0 is present and re-asserts the mwan3 ipset entry
# (which is ephemeral and lost on every reboot). When the flag is "0" or
# unset, fast no-op.
#
# Why this exists:
#   fw4's default input policy is DROP. Without a zone for tailscale0,
#   inbound VPN traffic hits no jump rule and is silently dropped. mwan3
#   marks reply packets for WAN egress unless the VPN CGNAT range is in its
#   connected-routes ipset — which is in-kernel only and cleared on every
#   boot, so it must be re-added each time. The matching CGI install flow
#   only adds these when the user enables the toggle in System Settings.
#
# Scope: Tailscale only. NetBird's parallel boot self-heal lives in
# qmanager_netbird_zone and is unrelated.
#
# Non-procd: runs once at boot, exits. No watch, no respawn.
#
# Usage:
#   /etc/init.d/qmanager_vpn_zone start|enable|disable
#
# Install location: /etc/init.d/qmanager_vpn_zone
# =============================================================================

START=99
STOP=10

. /usr/lib/qmanager/vpn_firewall.sh

# -----------------------------------------------------------------------------
# ensure_if_installed <binary> <zone_name> <device>
#   Calls vpn_fw_ensure_zone only when:
#     1. the VPN binary is on PATH (VPN is installed), AND
#     2. the fw4 zone is currently missing.
#   Both checks together make this a cheap idempotent no-op in the common case.
# -----------------------------------------------------------------------------
ensure_if_installed() {
    bin="$1"
    zone_name="$2"
    device="$3"

    if ! command -v "$bin" >/dev/null 2>&1; then
        logger -t qmanager_vpn_zone "$bin not installed, skipping $zone_name zone"
        return 0
    fi

    if vpn_fw_zone_exists "$zone_name"; then
        logger -t qmanager_vpn_zone "Zone '$zone_name' already present, no action"
        # Still re-assert mwan3 ipset (ephemeral, lost on reboot).
        vpn_fw_ensure_mwan3_exception
        return 0
    fi

    logger -t qmanager_vpn_zone "Zone '$zone_name' missing for installed $bin — creating"
    vpn_fw_ensure_zone "$zone_name" "$device"
}

start() {
    # Gate: only run when the user has explicitly enabled the workarounds.
    if [ "$(uci -q get quecmanager.tailscale_workarounds.enabled)" != "1" ]; then
        logger -t qmanager_vpn_zone "Tailscale workarounds disabled, skipping"
        return 0
    fi

    ensure_if_installed "tailscale" "tailscale" "tailscale0"
}

# stop() is intentionally a no-op. Zone removal on stop would force a
# firewall restart that breaks VPN reachability on every service reload,
# which is precisely the opposite of what this self-heal is for. Zone
# removal is owned by the disable path in system/settings.sh.
stop() {
    :
}
