# scripts/etc/init.d/qmanager_tower_failover
#!/bin/sh /etc/rc.common
# =============================================================================
# qmanager_tower_failover — Init Script for Tower Lock Failover Daemon
# =============================================================================
# Manages the tower lock failover daemon lifecycle. Non-procd: no automatic
# respawn, because the daemon intentionally exits after triggering failover
# or when locks are removed.
#
# All spawn/kill paths (CGI, cron, boot) converge through this script.
# The library functions tower_spawn_failover_watcher() and
# tower_kill_failover_watcher() call start/stop here.
#
# Usage:
#   /etc/init.d/qmanager_tower_failover start|stop|enable|disable
#
# Install location: /etc/init.d/qmanager_tower_failover
# =============================================================================

START=99
STOP=10

DAEMON="/usr/bin/qmanager_tower_failover"
PID_FILE="/tmp/qmanager_tower_failover.pid"
ACTIVATED_FLAG="/tmp/qmanager_tower_failover"
CONFIG="/etc/qmanager/tower_lock.json"

start() {
    # Guard: already running?
    if [ -f "$PID_FILE" ]; then
        local pid
        pid=$(cat "$PID_FILE" 2>/dev/null | tr -d ' \n\r')
        if [ -n "$pid" ] && kill -0 "$pid" 2>/dev/null; then
            return 0
        fi
        # Stale PID file
        rm -f "$PID_FILE"
    fi

    # Guard: daemon script exists and is executable?
    [ -x "$DAEMON" ] || return 1

    # Guard: config has failover enabled + active lock?
    if [ -f "$CONFIG" ]; then
        local fo lte nr
        fo=$(jq -r '.failover.enabled // false' "$CONFIG" 2>/dev/null)
        lte=$(jq -r '.lte.enabled // false' "$CONFIG" 2>/dev/null)
        nr=$(jq -r '.nr_sa.enabled // false' "$CONFIG" 2>/dev/null)
        if [ "$fo" != "true" ]; then
            return 0
        fi
        if [ "$lte" != "true" ] && [ "$nr" != "true" ]; then
            return 0
        fi
    else
        return 0
    fi

    # Clear previous activation flag
    rm -f "$ACTIVATED_FLAG"

    # Spawn daemon (double-fork to detach from any parent process group)
    ( "$DAEMON" </dev/null >/dev/null 2>&1 & )
    logger -t qmanager_tower_failover "Failover daemon started"
}

stop() {
    if [ -f "$PID_FILE" ]; then
        local pid
        pid=$(cat "$PID_FILE" 2>/dev/null | tr -d ' \n\r')
        if [ -n "$pid" ] && kill -0 "$pid" 2>/dev/null; then
            kill "$pid" 2>/dev/null
            logger -t qmanager_tower_failover "Failover daemon stopped (PID=$pid)"
        fi
        rm -f "$PID_FILE"
    fi
    rm -f "$ACTIVATED_FLAG"
}
