# scripts/etc/init.d/qmanager_wan_guard
#!/bin/sh /etc/rc.common
# =============================================================================
# qmanager_wan_guard — Init Script for WAN Interface Guard Daemon
# =============================================================================
# Launches the one-shot WAN guard daemon at boot. The daemon checks each
# rmnet WAN interface's profile number against modem-defined CIDs
# (AT+CGDCONT?) and disables any interface whose CID doesn't exist,
# preventing netifd from looping on phantom rmnet profiles.
#
# Non-procd: runs once and exits. No automatic respawn needed.
#
# Usage:
#   /etc/init.d/qmanager_wan_guard start|enable|disable
#
# Install location: /etc/init.d/qmanager_wan_guard
# =============================================================================

START=99
STOP=10

DAEMON="/usr/bin/qmanager_wan_guard"

start() {
    # Guard: daemon script exists and is executable
    [ -x "$DAEMON" ] || return 1

    # Spawn daemon (double-fork to detach from init process group)
    ( "$DAEMON" </dev/null >/dev/null 2>&1 & )
    logger -t qmanager_wan_guard "WAN guard daemon started"
}

stop() {
    # One-shot daemon — no persistent process to stop
    :
}
