# scripts/etc/init.d/qmanager
#!/bin/sh /etc/rc.common
# =============================================================================
# qmanager — OpenWRT Init Script (procd)
# =============================================================================
# Manages all QManager daemon lifecycles:
#   1. qmanager_ping     — Unified ping daemon (internet reachability + latency)
#   2. qmanager_poller   — Modem data poller (AT commands + system metrics)
#   3. qmanager_watchcat  — Connection health watchdog (future)
#
# Each daemon runs as a separate procd instance with independent respawn.
# Startup order matters: ping → poller → watchcat
#
# Usage:
#   /etc/init.d/qmanager start|stop|restart|enable|disable
#
# Install location: /etc/init.d/qmanager
# =============================================================================

START=99
STOP=10
USE_PROCD=1

start_service() {
    # Ensure required directories exist
    mkdir -p /var/lock
    mkdir -p /etc/qmanager
    mkdir -p /usr/lib/qmanager

    # Ensure all /usr/bin scripts are executable
    for f in /usr/bin/qmanager_*; do
        [ -f "$f" ] && chmod +x "$f"
    done

    # Ensure all CGI scripts are executable
    for f in /www/cgi-bin/quecmanager/*.sh /www/cgi-bin/quecmanager/*/*.sh; do
        [ -f "$f" ] && chmod +x "$f"
    done

    # Ensure auth config has restrictive permissions if it exists
    [ -f /etc/qmanager/auth.json ] && chmod 600 /etc/qmanager/auth.json

    # Create default long commands list if it doesn't exist
    if [ ! -f /etc/qmanager/long_commands.list ]; then
        cat > /etc/qmanager/long_commands.list << 'EOF'
# QManager Long Commands List
# Commands listed here will use extended timeouts and signal the poller to yield.
# One command keyword per line. Lines starting with # are comments.
QSCAN
QSCANFREQ
QFOTADL
EOF
    fi

    # Read log level from config (default: INFO)
    local log_level
    log_level=$(cat /etc/qmanager/log_level 2>/dev/null || echo "INFO")

    # -----------------------------------------------------------------
    # Instance 1: Ping Daemon (starts first — no dependencies)
    # -----------------------------------------------------------------
    procd_open_instance "ping"
    procd_set_param command /usr/bin/qmanager_ping
    procd_set_param env QLOG_LEVEL="$log_level"
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param respawn 3600 5 5
    procd_close_instance

    # -----------------------------------------------------------------
    # Instance 2: Poller Daemon (reads ping data for connectivity section)
    # -----------------------------------------------------------------
    procd_open_instance "poller"
    procd_set_param command /usr/bin/qmanager_poller
    procd_set_param env QLOG_LEVEL="$log_level"
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param respawn 3600 5 5
    procd_close_instance

    # -----------------------------------------------------------------
    # Instance 3: Watchcat Daemon (connection health monitor)
    # Only starts if enabled in UCI config.
    # -----------------------------------------------------------------
    local wc_enabled
    wc_enabled=$(uci -q get quecmanager.watchcat.enabled 2>/dev/null)
    if [ "$wc_enabled" = "1" ]; then
        procd_open_instance "watchcat"
        procd_set_param command /usr/bin/qmanager_watchcat
        procd_set_param env QLOG_LEVEL="$log_level"
        procd_set_param stdout 1
        procd_set_param stderr 1
        procd_set_param respawn 3600 5 5
        procd_close_instance
    fi

    # Tower failover daemon is managed by its own init.d script:
    # /etc/init.d/qmanager_tower_failover (enable it once after install)
}

stop_service() {
    # Clean up state files
    rm -f /tmp/qmanager_long_running
    rm -f /var/lock/qmanager.pid
    rm -f /tmp/qmanager_recovery_active
    rm -f /tmp/qmanager_watchcat.pid
    rm -f /tmp/qmanager_watchcat.lock
    # Note: Ping cache and history are left for post-mortem debugging
    # Note: Logs in /tmp/qmanager.log are preserved for post-mortem debugging
}

service_triggers() {
    procd_add_reload_trigger "qmanager"
}
