# scripts/etc/init.d/qmanager_dpi
#!/bin/sh /etc/rc.common
# qmanager_dpi — Video Optimizer (DPI evasion) service
# Manages nfqws lifecycle for video throttle bypass

START=99
STOP=10
USE_PROCD=1

. /usr/lib/qmanager/qlog.sh
. /usr/lib/qmanager/dpi_helper.sh

qlog_init "qmanager_dpi"

start_service() {
    local enabled
    enabled=$(uci -q get quecmanager.video_optimizer.enabled)
    [ "$enabled" = "1" ] || return 0

    # Check binary
    if ! dpi_check_binary; then
        qlog_error "nfqws binary not found at $NFQWS_BIN"
        return 1
    fi

    # Check kernel module
    if ! dpi_check_kmod; then
        qlog_error "Failed to load nfnetlink_queue kernel module"
        return 1
    fi

    # Check hostlist
    if [ ! -f "$DPI_HOSTLIST" ]; then
        qlog_error "Hostname list not found at $DPI_HOSTLIST"
        return 1
    fi

    # Read UCI config
    local quic_enabled iface
    quic_enabled=$(uci -q get quecmanager.video_optimizer.quic_enabled)
    iface=$(uci -q get quecmanager.video_optimizer.interface)
    DPI_INTERFACE="${iface:-wwan0}"

    # Remove stale rules, insert fresh
    dpi_remove_rules
    if ! dpi_insert_rules "$DPI_INTERFACE"; then
        qlog_error "Failed to insert nftables NFQUEUE rules"
        return 1
    fi

    # Build nfqws args
    local args="--qnum=$DPI_QUEUE_NUM"
    args="$args --hostlist=$DPI_HOSTLIST"
    args="$args --dpi-desync=split2"
    args="$args --dpi-desync-split-seqovl=1"
    args="$args --dpi-desync-split-pos=1"

    # Add QUIC desync if enabled (default: yes)
    if [ "$quic_enabled" != "0" ]; then
        args="$args --dpi-desync-udplen-increment=2"
    fi

    qlog_info "Starting nfqws on $DPI_INTERFACE (queue=$DPI_QUEUE_NUM)"

    procd_open_instance nfqws
    procd_set_param command $NFQWS_BIN $args
    procd_set_param respawn 3600 5 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance

    # Write PID file manually for uptime tracking (procd doesn't write one)
    sleep 1
    pidof nfqws > "$NFQWS_PID" 2>/dev/null || true

    qlog_info "nfqws service started"
}

stop_service() {
    qlog_info "Stopping nfqws service"
    dpi_remove_rules
    [ -f "$NFQWS_PID" ] && rm -f "$NFQWS_PID"
}

service_triggers() {
    procd_add_reload_trigger "quecmanager"
}
