# scripts/etc/init.d/qmanager_eth_link
#!/bin/sh /etc/rc.common
# =============================================================================
# qmanager_eth_link — Init Script for Ethernet Link Speed Limit Persistence
# =============================================================================
# Applies the configured ethernet link speed limit at boot using ethtool.
# Reads the limit from UCI (quecmanager.eth_link.speed_limit) and sets the
# advertise mode accordingly.
#
# Non-procd: runs once at boot and exits.
#
# Usage:
#   /etc/init.d/qmanager_eth_link start|enable|disable
#
# Install location: /etc/init.d/qmanager_eth_link
# =============================================================================

START=99
STOP=10

ETH_INTERFACE="eth0"
. /usr/lib/qmanager/ethtool_helper.sh

# Build hex advertise mask from ethtool's "Supported link modes" output.
# Maps mode names to bit positions (linux/ethtool.h), splits into hi/lo
# 32-bit halves so awk handles bit 47+ (2500baseT/Full) correctly.
start() {
    # Guard: ethtool must be available
    command -v ethtool >/dev/null 2>&1 || return 1

    # Guard: interface must exist
    [ -d "/sys/class/net/$ETH_INTERFACE" ] || return 1

    # Read configured speed limit from UCI
    limit=$(uci get quecmanager.eth_link.speed_limit 2>/dev/null)

    # Wait for interface to be ready
    sleep 2

    # Apply speed limit
    case "$limit" in
        "10")   ethtool -s "$ETH_INTERFACE" advertise 0x003 autoneg on 2>/dev/null ;;
        "100")  ethtool -s "$ETH_INTERFACE" advertise 0x00f autoneg on 2>/dev/null ;;
        "1000") ethtool -s "$ETH_INTERFACE" advertise 0x02f autoneg on 2>/dev/null ;;
        *)
            advertise=$(get_supported_advertise_hex)
            if [ -n "$advertise" ]; then
                ethtool -s "$ETH_INTERFACE" advertise "$advertise" autoneg on 2>/dev/null
            else
                ethtool -s "$ETH_INTERFACE" autoneg on 2>/dev/null
            fi
            ;;
    esac

    ethtool -r "$ETH_INTERFACE" 2>/dev/null

    logger -t qmanager_eth_link "Applied speed limit: ${limit:-auto}"
}

stop() {
    # Reset to auto on stop — advertise all supported modes
    if command -v ethtool >/dev/null 2>&1; then
        advertise=$(get_supported_advertise_hex)
        if [ -n "$advertise" ]; then
            ethtool -s "$ETH_INTERFACE" advertise "$advertise" autoneg on 2>/dev/null
        else
            ethtool -s "$ETH_INTERFACE" autoneg on 2>/dev/null
        fi
        ethtool -r "$ETH_INTERFACE" 2>/dev/null
        logger -t qmanager_eth_link "Reset to auto speed"
    fi
}
